# yuga-check.service
[Unit]
Description=Yuga Check Tool

[Service]
Type=oneshot
ExecStart=/usr/sbin/yuga-check

[Install]
WantedBy=default.target